# homestead.root/etc/cron.hourly/homestead-log-cleanup
#!/bin/bash
WEEK_COUNT=168
ONE_GIG=1073741824
max_log_directory_size=$ONE_GIG
. /etc/clearwater/config

python2.7 /usr/share/clearwater/bin/log_cleanup.py /var/log/homestead --prefix ramtrace_ --maxsize $max_log_directory_size --count 2
python2.7 /usr/share/clearwater/bin/log_cleanup.py /var/log/homestead --prefix homestead_ --maxsize $max_log_directory_size --count $WEEK_COUNT
python2.7 /usr/share/clearwater/bin/log_cleanup.py /var/log/homestead --prefix access_ --maxsize $max_log_directory_size --count $WEEK_COUNT
python2.7 /usr/share/clearwater/bin/log_cleanup.py /var/lib/cassandra --suffix hprof --count 2
